# data/siglo.service
[Unit]
Description=siglo service
[Service]
ExecStart=siglo --daemon
