# data/siglo.service
[Unit]
Description=siglo service
[Service]
ExecStart=siglo --start
ExecStop=siglo --stop
Environment=PYTHONUNBUFFERED=1
